# supervisor/bkp_mydumper.conf
[program:bkp_mydumper]
command=/home/noc/bkp_mydumper/run.py
priority=1
user=root
redirect_stderr=true
autostart=true
autorestart=true
startsecs = 1
